# ansible/roles/wallpaper/templates/rotate-wallpaper.service.j2
[Unit]
Description=Rotate desktop wallpaper

[Service]
Type=oneshot
ExecStart={{ home }}/.local/bin/rotate-wallpaper.sh
