# roles/wallpaper/templates/rotate-wallpaper.timer.j2
[Unit]
Description=Rotate wallpaper timer

[Timer]
OnBootSec=1m
OnUnitActiveSec={{ wallpaper_interval }}

[Install]
WantedBy=timers.target
